# service_files/arch/atxpi-gpio-shutdown.service
[Unit]
Description=atxpi gpio shutdown script

[Service]
Type=simple
ExecStart=/opt/atxpi/shutdown.sh

[Install]
WantedBy=multi-user.target
